# scrapper.service
[Unit]
Description=Price Scrapper
Wants=network-online.target
After=network-online.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/python-scrapper
EnvironmentFile=/var/www/python-scrapper/.env
Environment="PATH=/var/www/python-scrapper/venv/bin"
ExecStart=/var/www/python-scrapper/venv/bin/gunicorn -w 3 -k gthread -b 0.0.0.0:${PORT} --access-logfile - --error-logfile - app.url_manager:app
Restart=always
RestartSec=5
TimeoutStopSec=15

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target